# 40-mouse-razer.rules
ACTION=="add", ATTRS{product}=="Razer Naga", SYMLINK+="NagaRazerMouse", RUN+="/usr/local/bin/udev-naga.sh"

